# .config/supervisord/gunicorn.conf
[program:gunicorn]
command=gunicorn -c /home/docker/.config/config/gunicorn.py app:app