# .config/supervisord/nginx.conf
[program:nginx]
command=nginx -c /home/docker/.config/config/nginx.conf